# systemd/user/picframe.service
# PicFrame Display Service
# Location: ~/.config/systemd/user/picframe.service
#
# Install with:
#   mkdir -p ~/.config/systemd/user
#   cp picframe.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable picframe.service

[Unit]
Description=PictureFrame on Pi

[Service]
# Redirect Xorg log to /dev/null to prevent log bloat
# The GPU driver logs EDID info ~15 times/sec causing massive log growth
# Must specify full Xorg path for xinit to pass server args correctly
ExecStart=xinit /home/pi/picframe_3.0/app_control/start_picframe_app.sh -- /usr/lib/xorg/Xorg -logfile /dev/null
Restart=always

[Install]
WantedBy=default.target
